# crontab
*/5 * * * * php /var/www/bin/console app:run 2>&1
